# meshcore-repeater.service
[Unit]
Description=MeshCore Repeater service
Documentation=man:systemd.service(5)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/meshcore-repeater --service repeater /etc/meshcore-netcore/config.toml
WorkingDirectory=/var/lib/meshcore
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3
TimeoutStopSec=30
KillMode=control-group
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meshcore-repeater
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1
Environment=MESHCORE_LOG_DIR=/var/log/meshcore-netcore
Environment=MESHCORE_LOG_FILE=meshcore-repeater.log
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=yes
ReadWritePaths=/etc/meshcore-netcore /var/lib/meshcore /var/log/meshcore-netcore

[Install]
WantedBy=multi-user.target
